# hn-summaries.service
[Unit]
Description=HN Summaries — AI-summarized Hacker News RSS feed
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/app
# Absolute bun path so PATH resolution is irrelevant under systemd.
# Bun runs the TS directly (no tsx) and provides Bun.WebView for browser extraction.
ExecStart=/home/exedev/.bun/bin/bun /home/exedev/app/index.ts
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=6
Environment=NODE_ENV=production
Environment=PORT=8000
Environment=PUBLIC_URL=https://hn.rlew.io
# Chrome/Chromium for Bun.WebView's browser-extraction tier (Linux drives it over CDP).
# The app resolves the binary at startup (system browser, else newest Playwright cache),
# so no version-pinned path is hardcoded here. Install one rootless with
# `bun run install-browser`. Uncomment to force a specific binary:
# Environment=BUN_CHROME_PATH=/path/to/chrome
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hn-summaries

[Install]
WantedBy=multi-user.target
